# systemd/hotspot-shutdown.service
[Unit]
Description=Stop Hotspot after timeout

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop setup-hotspot.target
