# systemd/hotspot-shutdown.timer
[Unit]
Description=Shutdown hotspot after 15 minutes

[Timer]
OnActiveSec=15min
Unit=hotspot-shutdown.service

[Install]
WantedBy=multi-user.target
